# services/mcpd.service
[Unit]
Description=MCP Daemon - Centralized MCP tool registry
Documentation=https://github.com/youruser/mcpd
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mcpd start
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/mcpd
PrivateTmp=true

[Install]
WantedBy=default.target
